# vsftpd/super-vsftpd.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command= service vsftpd start
autorestart= true