# rsync/conf/super-rsync.conf
[supervisord]
nodaemon=true

[program:rsync]
command= /usr/bin/rsync --no-detach --daemon --config /etc/rsyncd.conf "$@"
autorestart= true