# apache/super-apache.conf
[supervisord]
nodaemon=true

[program:apache2]
command= service apache2 start
autorestart= true